# delete_old_amavis_mail.crontab
@daily amavis /usr/local/sbin/delete_old_amavis_mail
